# init/deb/DEBIAN/prerm
#!/bin/bash

systemctl stop remote-admin
systemctl disable remote-admin

rm -f /opt/remote-admin/*.*
